# server/script/trait-o-matic.in
#!/bin/bash
# 
# chkconfig: 2345 20 80 
# description: Trait-o-matic service

export USER="%USER%"
export LOG="%LOG%"
export CORE="%CORE%"
export DATA="%DATA%"
SERVER="$CORE/trait-o-matic-server.py"
SERVER_PID="%TMP%/server.pid"

start() {
    echo Starting Trait-o-matic server.
    start-stop-daemon --pidfile $SERVER_PID --user $USER --exec `which python` --background --chuid $USER --startas $SERVER --start -- --pidfile=$SERVER_PID --stderr=$LOG/server.log
}

stop() {
    echo Stopping Trait-o-matic server.
    start-stop-daemon --pidfile $SERVER_PID --user $USER --exec `which python` --stop --retry 2
}

case "$1" in
    start)
	start
	;;
    stop)
	stop
	;;
    restart)
	stop
	start
	;;
    *)
	echo "Usage: $0 start|stop|restart"
	exit 1
	;;
esac
exit $?
